# Ansible_project/roles/mcdn-agent-py/files/mcdn-agent-py.service
[Unit]
Description=MCDN Agent Python Service
After=network.target docker.service

[Service]
Type=simple
ExecStart={service_dir}/startup.sh
WorkingDirectory={service_dir}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
